# docker/s6/s6-rc.d/snapotter/run
#!/command/with-contenv sh
cd /app
exec s6-setuidgid snapotter pnpm --filter @snapotter/api run start
