# docker/s6/s6-rc.d/snapotter/run
#!/command/with-contenv sh
cd /app/apps/api
exec s6-setuidgid snapotter ./node_modules/.bin/tsx --import ./src/tracing.ts --import ./src/instrument.ts src/index.ts
