# docker/s6/s6-rc.d/postgres-init/up
#!/command/with-contenv sh
/usr/local/bin/embedded-postgres-bootstrap.sh
